# deployment/systemd_service.service
[Unit]
Description=ONDC BAP FastAPI Application
After=network.target

[Service]
Type=exec
User=www-data
Group=www-data
WorkingDirectory=/opt/ondc-bap
Environment=PATH=/opt/ondc-bap/.venv/bin
ExecStart=/opt/ondc-bap/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 